# bin/zfor
#!/bin/sh
#
# zfor	Control script for Z Fail-over Resolver
#
# chkconfig: - 85 15
# description: Z Fail-over Resolver

start()
{
	sudo -u nobody /usr/local/bin/zfor-start
}

stop()
{
	sudo -u nobody /usr/local/bin/zfor-stop
}

dump()
{
	/usr/local/bin/zfor-conf
}

stat()
{
	/usr/local/bin/zfor-stat
}

case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	restart)
		stop
		sleep 1
		start
		;;
	dump)
		dump
		;;
	stat)
		stat
		;;
	*)
		echo "Usage: $0 {start|stop|restart|dump|stat}"
esac

exit 0

